# home/chandradeep/.config/autostart/ibus-daemon.desktop
[Desktop Entry]
Exec=/usr/bin/ibus-daemon -drx
Icon=
Name=ibus-daemon
Path=
Terminal=False
Type=Application
